# resources/build/etc/systemd/system/fb-miniserver.service
[Unit]
Description=FastyBird MiniServer
After=multi-user.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=miniserver
Group=miniserver

ExecStart=/usr/bin/php8.0 /usr/lib/miniserver/vendor/bin/fb-console fb:web-server:start
ExecStop=/bin/kill -INT $MAINPID
ExecReload=/bin/kill -TERM $MAINPID

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
