# resources/build/DEBIAN/postinst
#!/bin/sh
#     Copyright 2021. FastyBird s.r.o.
#
#     Licensed under the Apache License, Version 2.0 (the "License");
#     you may not use this file except in compliance with the License.
#     You may obtain a copy of the License at
#
#         http://www.apache.org/licenses/LICENSE-2.0
#
#     Unless required by applicable law or agreed to in writing, software
#     distributed under the License is distributed on an "AS IS" BASIS,
#     WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
#     See the License for the specific language governing permissions and
#     limitations under the License.

set -e

echo "Finalizing configuration"
sudo usermod -a -G miniserver_manager miniserver_manager
sudo chown miniserver_manager:miniserver_manager /etc/miniserver-manager/ -R
sudo chown miniserver_manager:miniserver_manager /usr/lib/miniserver-manager/ -R
sudo chown miniserver_manager:miniserver_manager /var/log/miniserver-manager/ -R
sudo chown miniserver_manager:miniserver_manager /var/tmp/miniserver-manager/ -R
sudo chown miniserver_manager:miniserver_manager /var/lib/miniserver_manager/ -R
sudo ln -s /usr/lib/miniserver-manager/vendor/fastybird/bootstrap/bin/fb-console /usr/bin/fb-miniserver-manager
echo "Installation completed"

echo "Enabling daemon..."
sudo pidof systemd && sudo systemctl enable fb-miniserver-manager || echo "Systemctl not found"

echo "Daemon starting..."
sudo pidof systemd && sudo systemctl start fb-miniserver-manager || echo
echo -e "\e[96mFastyBird MiniServer Gateway \e[92mhas been installed. Have a nice day \e[93m\e[5m:)\e[25m\e[39m"
